# systemd/llama-server.service
[Unit]
Description=llama.cpp server (CPU-only, AVX-512, MoE)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/opt/llama/llama.env
Environment=LD_LIBRARY_PATH=/opt/llama/lib

# threads (geracao) < threads-batch (prompt) de proposito: medido nesta maquina,
# a geracao e limitada por banda de memoria e satura com ~10 threads (11.37 tok/s
# contra 10.95 com 20), enquanto o prompt e compute-bound e escala com cores
# (80.4 tok/s com 20 contra 51.3 com 10). Separar os dois pega o melhor de cada.
#
# numactl prende a UM socket: medido, usar os 2 sockets DERRUBA a geracao pela
# metade (4.67 vs 10.87 tok/s) por custo de sincronizacao via UPI.
ExecStart=/usr/bin/numactl --cpunodebind=${NUMA_NODE} --membind=${NUMA_NODE} \
    /opt/llama/bin/llama-server \
    --model ${MODEL} \
    --threads ${THREADS} \
    --threads-batch ${THREADS_BATCH} \
    --ctx-size ${CTX_SIZE} \
    --numa numactl \
    --no-mmap \
    --jinja \
    --host 127.0.0.1 --port 8081

WorkingDirectory=/opt/llama
Restart=on-failure
RestartSec=10
TimeoutStartSec=900
LimitNOFILE=65536
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
