# systemd/llama-server.service
[Unit]
Description=llama.cpp server (CPU-only, AVX-512, MoE)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/opt/llama/llama.env
Environment=LD_LIBRARY_PATH=/opt/llama/lib

# threads (geracao) < threads-batch (prompt) de proposito: medido nesta maquina,
# a geracao e limitada por banda de memoria e satura com ~10 threads (11.37 tok/s
# contra 10.95 com 20), enquanto o prompt e compute-bound e escala com cores
# (80.4 tok/s com 20 contra 51.3 com 10). Separar os dois pega o melhor de cada.
#
# SOBRE NUMA -- correcao de uma conclusao anterior errada: este container LXC so
# enxerga os CPUs 0-19, todos do node 0 (Cpus_allowed_list: 0-19, nproc=20), e
# 'numactl --cpunodebind=1' falha aqui com "sched_setaffinity: Invalid argument".
# O teste antigo de "2 sockets, 40 threads" que rendeu 4.67 tok/s NAO mediu dois
# sockets: mediu 40 threads disputando 20 CPUs, ou seja oversubscription de 2x.
# Prender a um socket continua certo, mas por ser o unico disponivel ao
# container, e nao por custo de sincronizacao via UPI. Para usar de fato os dois
# sockets seria preciso reconfigurar os cores do LXC no host Proxmox.
#
# Todos os argumentos ficam explicitos abaixo. NAO reintroduzir uma variavel
# tipo ${EXTRA_ARGS} com varios argumentos dentro: o systemd nao faz word
# splitting em ${VAR}, entao a string inteira chega como UM argumento so e o
# numactl morre com "invalid argument: --jinja --cache-reuse 256".
ExecStart=/usr/bin/numactl --cpunodebind=${NUMA_NODE} --membind=${NUMA_NODE} \
    /opt/llama/bin/llama-server \
    --model ${MODEL} \
    --threads ${THREADS} \
    --threads-batch ${THREADS_BATCH} \
    --ctx-size ${CTX_SIZE} \
    --parallel ${N_PARALLEL} \
    --numa numactl \
    --no-mmap \
    --jinja \
    --host 127.0.0.1 --port 8081

WorkingDirectory=/opt/llama
Restart=on-failure
RestartSec=10
TimeoutStartSec=900
LimitNOFILE=65536
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
